# templates/etc/init.d/application-initscript.erb
#!/bin/bash
#
# tomcat7 This shell script takes care of starting and stopping Tomcat
#
# chkconfig: - 80 20
#
### BEGIN INIT INFO
# Provides: tomcat7
# Required-Start: $network $syslog
# Required-Stop: $network $syslog
# Default-Start:
# Default-Stop:
# Description: Release implementation for Servlet 2.5 and JSP 2.1
# Short-Description: start and stop tomcat
### END INIT INFO

. /etc/sysconfig/<%= @application_name %>

case $1 in
    start)
           su -c "<%= @catalina_home %>/bin/catalina.sh start" $CATALINA_USER
        ;;
    stop)
           su -c "<%= @catalina_home %>/bin/catalina.sh stop -force" $CATALINA_USER
        ;;
    restart)
           su -c "<%= @catalina_home %>/bin/catalina.sh stop -force" $CATALINA_USER
           su -c "<%= @catalina_home %>/bin/catalina.sh start" $CATALINA_USER
        ;;
    status)
          if [ -f $CATALINA_PID ]; then
              PID=`cat $CATALINA_PID`
              if [ -z "`ps axf | grep ${PID} | grep -v grep`" ]; then
                  exit 1
                  printf "%s\n" "Process dead but pidfile exists"
              else
                  echo "Running"
              fi
          else
              printf "%s\n" "Service not running"
              exit 3
          fi
        ;;
esac